# scripts/etc/systemd/system/linker-upgrade.service
[Unit]
Description=LinkerBot Program Upgrade System - Web-based application upgrade service
Documentation=https://github.com/linker-bot/linker-upgrader
After=network.target network-online.target
Wants=network-online.target
RequiresMountsFor=/opt

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/linker-upgrader
ExecStart=/opt/linker-upgrader/linker-upgrader -config /etc/linker-upgrader/config.json
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=10s

# 安全设置
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 环境变量
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=/root"

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=linker-upgrader

# 进程管理
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
Alias=linker-upgrader.service
